# deploy/ai-agent.service
[Unit]
Description=Telegram AI coding agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/ai-coding-agent
Environment=HOME=/root
EnvironmentFile=/etc/ai-agent/ai-agent.env
ExecStart=/opt/ai_coding_venv/bin/python /opt/ai-coding-agent/agent.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
